# files/fix-nhc.service
[Unit]
Description=Fix nhc.conf
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fix-nhc.py
TimeoutSec=60
