# files/fix-nhc.service
[Unit]
Description=Fix nhc.conf
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/fix-nhc.py /etc/nhc/nhc.conf
TimeoutSec=60

[Install]
WantedBy=multi-user.target
